# supervisord.conf
[supervisord]
nodaemon=true

[program:daemon-coin]
directory=/usr/local/bin
command=SperoCoind -datadir=/coin/data -conf=/coin/sperocoin.conf -txindex
autostart=true
autorestart=true