# ubuntu-1604-dbd48/supervisord.conf
[supervisord]
nodaemon=true

[program:daemon-coin]
directory=/usr/local/bin
command=SperoCoind -datadir=/coin/data -conf=/coin/data/sperocoin.conf
autostart=true
autorestart=true